# etc/cron.d/mastodon
PATH="/bin:/usr/bin:/sbin:/usr/sbin:/usr/local/bin:/home/mastodon/.rbenv/shims/"
RAILS_ENV=production
NUM_DAYS=5
#
# Clean up
#   https://docs.joinmastodon.org/admin/setup/
0 6 * * *  mastodon (echo "____________________________" >> /var/log/mastodon.log; date >>/var/log/mastodon.log 2>&1 ; cd /home/mastodon/live && bin/tootctl media remove >>/var/log/mastodon.log 2>&1)
0 12 * * 7 mastodon (echo "____________________________" >> /var/log/mastodon.log; date >>/var/log/mastodon.log 2>&1 ; cd /home/mastodon/live && bin/tootctl preview_cards remove >>/var/log/mastodon.log 2>&1)
